# files/audio/.config/systemd/user/darkice.service
[Unit]
Description=Darkice streams alsa monitor

[Service]
ExecStart=/usr/bin/darkice -c /home/nemo/.darkice.cfg
