# files/audio/.config/systemd/user/icecast.service
[Unit]
Description=Streams whatever darkice gets
Requires=darkice.service

[Service]
ExecStart=/usr/bin/icecast -c /home/nemo/.icecast.xml
